# home/.config/systemd/user/update-cursor.service
[Unit]
Description=Update Cursor

[Service]
ExecStart=/home/vitor/build/cursor/update-cursor.sh
Type=oneshot

[Install]
WantedBy=default.target
